# conf/full_raspberrypi_bcm27xx_bcm2712/files/etc/uci-defaults/93_osi_deploy_guard_init
#!/bin/sh
set -eu

# This is a ROM-owned, fresh-image initializer.  Test roots are accepted only
# by the explicit repair-program gate; production always uses fixed mounts and
# the trusted system Node binary.
if [ "${OSI_REPAIR_PROGRAM_MODE:-0}" = 1 ] && [ "${OSI_DEPLOY_ARTIFACT_MODE:-}" = test ]; then
    ROM_ROOT=${OSI_ROM_ROOT:-/rom}
    DATA_ROOT=${OSI_DATA_ROOT:-/data}
    RUN_ROOT=${OSI_RUN_ROOT:-/var/run}
    NODE=${OSI_NODE:-node}
    REBOOT_COMMAND=${OSI_REBOOT_COMMAND:-/sbin/reboot}
else
    ROM_ROOT=/rom
    DATA_ROOT=/data
    RUN_ROOT=/var/run
    NODE=/usr/bin/node
    REBOOT_COMMAND=/sbin/reboot
fi
PATH=/usr/sbin:/usr/bin:/sbin:/bin
export PATH
PROVENANCE="$ROM_ROOT/usr/share/osi-deploy/factory-image-provenance.json"
MANIFEST="$ROM_ROOT/usr/share/osi-deploy/image-guard-manifest.json"
PROVENANCE_CODEC="$ROM_ROOT/usr/libexec/osi-factory-image-provenance.js"
PROVENANCE_CLI="$ROM_ROOT/usr/libexec/osi-factory-image-provenance-cli.js"
STATE_CLI="$ROM_ROOT/usr/libexec/osi-deployment-state-cli.js"
SEED_HELPER="$ROM_ROOT/usr/libexec/osi-factory-database-seed-cli.js"
ACK_AUDIT_CLI="$ROM_ROOT/usr/libexec/osi-audit-command-ack-state.js"
PROTOCOL_CLI="$ROM_ROOT/usr/libexec/osi-sync-protocol-capability-cli.js"
RESULT="$RUN_ROOT/osi-factory-image-verification.json"
STATE_ROOT="$DATA_ROOT/osi-deploy"
DB="$DATA_ROOT/db/farming.db"
STATE="$STATE_ROOT/deployment-state.json"
LINEAGE="$STATE_ROOT/factory-database-lineage.json"

# The image-baseline operation ID is deterministic from the immutable image
# build identity.  The future state verb derives the same ID; callers never
# supply an arbitrary operation ID.

fail() { logger -t osi-deploy-guard-init "$1" 2>/dev/null || :; return 1; }
safe_path() {
    candidate=$1
    case "$candidate" in
        /*) ;;
        *) fail "path must be absolute: $candidate" ;;
    esac
    while [ "$candidate" != / ]; do
        [ ! -L "$candidate" ] || fail "symlink path component: $candidate"
        candidate=$(dirname "$candidate")
    done
}
safe_path "$ROM_ROOT"; safe_path "$DATA_ROOT"; safe_path "$RUN_ROOT"
for candidate in "$PROVENANCE" "$MANIFEST" "$PROVENANCE_CODEC" "$PROVENANCE_CLI" "$STATE_CLI" "$SEED_HELPER" "$ACK_AUDIT_CLI" "$PROTOCOL_CLI" "$STATE_ROOT" "$DB"; do safe_path "$candidate"; done
safe_path "$0"
[ -f "$0" ] && cmp -s "$0" "$ROM_ROOT/etc/uci-defaults/93_osi_deploy_guard_init" || fail "overlay 93 differs from ROM authority"
[ -d "$DATA_ROOT" ] && [ ! -L "$DATA_ROOT" ] || fail "persistent data root missing or unsafe"
[ -f "$PROVENANCE" ] && [ ! -L "$PROVENANCE" ] || fail "factory provenance missing or unsafe"
[ -f "$MANIFEST" ] && [ ! -L "$MANIFEST" ] || fail "image guard manifest missing or unsafe"
[ -f "$PROVENANCE_CODEC" ] && [ ! -L "$PROVENANCE_CODEC" ] || fail "provenance codec missing or unsafe"
for resident in "$PROVENANCE_CLI" "$STATE_CLI" "$SEED_HELPER" "$ACK_AUDIT_CLI" "$PROTOCOL_CLI"; do
    [ -f "$resident" ] && [ ! -L "$resident" ] || fail "resident authority missing or unsafe: $resident"
done
[ ! -e "$STATE_ROOT" ] || fail "deployment state already exists; fresh-image initializer refused"
for sidecar in "$DB-wal" "$DB-shm" "$DB-journal"; do
    safe_path "$sidecar"
    [ ! -e "$sidecar" ] || fail "database sidecar already exists"
done
[ ! -e "$DB" ] || fail "database already exists; fresh-image initializer refused"

if [ "${OSI_REPAIR_PROGRAM_MODE:-0}" = 1 ] && [ "${OSI_DEPLOY_ARTIFACT_MODE:-}" = test ] && [ -n "${OSI_IMAGE_PROFILE:-}" ]; then
    PROFILE=$OSI_IMAGE_PROFILE
else
    PROFILE=$($NODE -e 'const fs=require("fs"); const p=JSON.parse(fs.readFileSync(process.argv[1])); process.stdout.write(p.profile);' "$PROVENANCE")
fi
"$NODE" "$ROM_ROOT/usr/libexec/osi-factory-image-provenance-cli.js" verify-runtime \
    --factory-provenance "$PROVENANCE" \
    --image-guard-manifest "$MANIFEST" \
    --initializer "$ROM_ROOT/etc/uci-defaults/93_osi_deploy_guard_init" \
    --factory-seed "$ROM_ROOT/usr/share/db/farming.db" \
    --factory-seed-helper "$ROM_ROOT/usr/libexec/osi-factory-database-seed-cli.js" \
    --ack-audit-cli "$ROM_ROOT/usr/libexec/osi-audit-command-ack-state.js" \
    --protocol-cli "$ROM_ROOT/usr/libexec/osi-sync-protocol-capability-cli.js" \
    --expected-profile "$PROFILE" \
    --result "$RESULT"
VERIFICATION_SHA256=$($NODE -e 'const fs=require("fs"); const c=require("crypto"); process.stdout.write(c.createHash("sha256").update(fs.readFileSync(process.argv[1])).digest("hex"));' "$RESULT")
VERIFICATION_NONCE=$($NODE -e 'const fs=require("fs"); const p=JSON.parse(fs.readFileSync(process.argv[1])); if (!/^[0-9a-f]{64}$/.test(p.nonce)) process.exit(1); process.stdout.write(p.nonce);' "$RESULT")

BASELINE_ID=$($NODE - "$PROVENANCE" "$MANIFEST" "$PROVENANCE_CODEC" <<'NODE'
const fs = require('node:fs');
const [provenancePath, manifestPath, codecPath] = process.argv.slice(2);
const codec = require(codecPath);
const p = codec.readCanonicalJson(provenancePath, 'factory provenance');
const manifest = fs.readFileSync(manifestPath);
const parsedManifest = codec.readCanonicalJson(manifestPath, 'image-guard manifest');
codec.assertProfileRelation(p, p.profile);
codec.validateManifest(parsedManifest, p.profile);
if (p.imageBuildId !== parsedManifest.imageBuildId || p.imageGuardManifestSha256 !== codec.hashFile(manifestPath)) process.exit(1);
process.stdout.write(codec.factoryBaselineId(p.profile, manifest));
NODE
)
RECEIPT="$STATE_ROOT/receipts/$BASELINE_ID.factory-seed.json"
safe_path "$STATE"; safe_path "$RECEIPT"; safe_path "$LINEAGE"

# Commit-1's state CLI deliberately rejects this future verb.  This checkpoint
# therefore fails closed before any database/state mutation; the full runtime
# adds the image-baseline verb and owns seed receipts and lineage CAS.
"$NODE" "$ROM_ROOT/usr/libexec/osi-deployment-state-cli.js" initialize-image-baseline \
    --root "$STATE_ROOT" \
    --state "$STATE" \
    --receipts "$STATE_ROOT/receipts" \
    --guard-marker "$STATE_ROOT/guard-installed.json" \
    --image-manifest "$MANIFEST" \
    --factory-verification-result "$RESULT" \
    --factory-seed "$ROM_ROOT/usr/share/db/farming.db" \
    --factory-seed-helper "$ROM_ROOT/usr/libexec/osi-factory-database-seed-cli.js" \
    --factory-seed-receipt-out "$RECEIPT" \
    --database-lineage-out "$LINEAGE" \
    --expected-verification-sha256 "$VERIFICATION_SHA256" \
    --expected-verification-nonce "$VERIFICATION_NONCE" \
    --expected-profile "$PROFILE"

"$REBOOT_COMMAND"
