# conf/full_raspberrypi_bcm27xx_bcm2712/files/etc/uci-defaults/95_osi_bootstrap_enable
#!/bin/sh
/etc/init.d/osi-bootstrap enable
