# conf/full_raspberrypi_bcm27xx_bcm2712/files/etc/uci-defaults/96_osi_server_config
#!/bin/sh
# Only run if not already configured
[ -n "$(uci -q get osi-server.cloud.device_eui 2>/dev/null)" ] && exit 0

# Derive device EUI from Ethernet MAC
ETH_MAC=$(ip link show eth0 2>/dev/null | awk '/ether/{print $2}' | tr -d ':' | tr 'a-f' 'A-F')
[ -z "$ETH_MAC" ] && ETH_MAC=$(ip link show wlan0 2>/dev/null | awk '/ether/{print $2}' | tr -d ':' | tr 'a-f' 'A-F')

uci -q batch <<EOF
set osi-server.cloud=osi_server
set osi-server.cloud.enabled=0
set osi-server.cloud.device_eui=$ETH_MAC
set osi-server.cloud.device_type=GATEWAY
set osi-server.cloud.firmware_version=0.5.0
set osi-server.cloud.server_host=
set osi-server.cloud.mqtt_password=
set osi-server.cloud.openagri_weather_url=
set osi-server.cloud.openagri_weather_username=
set osi-server.cloud.openagri_weather_password=
set osi-server.cloud.openagri_weather_bearer_token=
set osi-server.cloud.openagri_weather_radius_km=10
set osi-server.cloud.openagri_weather_current_cache_minutes=30
set osi-server.cloud.openagri_weather_forecast_cache_minutes=120
commit osi-server
EOF
exit 0
