# conf/full_raspberrypi_bcm27xx_bcm2712/files/etc/uci-defaults/97_osi_db_seed
#!/bin/sh
set -eu

# 97 is a completion/verification hook, never an authority creator.  The
# only legal authority is the verified 93 image-baseline chain; test roots and
# adapters are available solely inside the explicit repair-program gate.
if [ "${OSI_REPAIR_PROGRAM_MODE:-0}" = 1 ] && [ "${OSI_DEPLOY_ARTIFACT_MODE:-}" = test ]; then
    ROM_ROOT=${OSI_ROM_ROOT:-/rom}
    DATA_ROOT=${OSI_DATA_ROOT:-/data}
    NODE=${OSI_NODE:-node}
else
    ROM_ROOT=/rom
    DATA_ROOT=/data
    NODE=/usr/bin/node
fi
DB="$DATA_ROOT/db/farming.db"
STATE_ROOT="$DATA_ROOT/osi-deploy"
SEED="$ROM_ROOT/usr/share/db/farming.db"
PROVENANCE="$ROM_ROOT/usr/share/osi-deploy/factory-image-provenance.json"
STATE="$STATE_ROOT/deployment-state.json"
RECEIPT="$STATE_ROOT/receipts/factory-seed.json"
LINEAGE="$STATE_ROOT/database-lineage.json"
PUBLICATION="$LINEAGE.seed-publication.json"

fail() { logger -t osi-db-seed "$1" 2>/dev/null || :; return 1; }
safe_path() {
    candidate=$1
    case "$candidate" in
        /*) ;;
        *) fail "path must be absolute: $candidate" ;;
    esac
    while [ "$candidate" != / ]; do
        [ ! -L "$candidate" ] || fail "symlink path component: $candidate"
        candidate=$(dirname "$candidate")
    done
}
for candidate in "$ROM_ROOT" "$DATA_ROOT" "$DB" "$STATE_ROOT" "$STATE" "$RECEIPT" "$LINEAGE" "$PUBLICATION" "$SEED" "$PROVENANCE"; do safe_path "$candidate"; done
[ -d "$DATA_ROOT" ] && [ ! -L "$DATA_ROOT" ] || fail "persistent data root missing or unsafe"
[ -f "$SEED" ] && [ ! -L "$SEED" ] || fail "ROM seed missing or unsafe"
[ -f "$PROVENANCE" ] && [ ! -L "$PROVENANCE" ] || fail "ROM provenance missing or unsafe"

# Sidecars are checked before the existing-main-file branch.  An abandoned
# WAL/SHM/journal is never made harmless by an apparently valid main file.
for sidecar in "$DB-wal" "$DB-shm" "$DB-journal"; do
    safe_path "$sidecar"
    [ ! -e "$sidecar" ] || fail "sidecar exists; refusing seed"
done

# A foreign/mismatched state chain is not idempotence.  Validate all durable
# evidence before accepting an existing database (or attempting any seed).
"$NODE" - "$STATE" "$RECEIPT" "$LINEAGE" "$DB" "$PROVENANCE" "$SEED" "$PUBLICATION" "$ROM_ROOT/usr/libexec/osi-factory-database-seed.js" "$ROM_ROOT/usr/libexec/osi-deployment-state.js" <<'NODE'
const fs = require('node:fs');
const crypto = require('node:crypto');
const [statePath, receiptPath, lineagePath, databasePath, provenancePath, seedPath, publicationPath, seedLibraryPath, stateLibraryPath] = process.argv.slice(2);
const sha = value => crypto.createHash('sha256').update(value).digest('hex');
const canonical = value => Array.isArray(value) ? `[${value.map(canonical).join(',')}]`
  : (value && typeof value === 'object' ? `{${Object.keys(value).sort().map(key => `${JSON.stringify(key)}:${canonical(value[key])}`).join(',')}}` : JSON.stringify(value));
function regular(file, label, expectedMode = 0o600) {
  const stat = fs.lstatSync(file);
  if (!stat.isFile() || stat.isSymbolicLink()) throw new Error(`${label} must be a regular nonsymlink file`);
  if (expectedMode !== null && stat.uid !== process.getuid()) throw new Error(`${label} must be owned by the runtime user`);
  if (expectedMode !== null && (stat.mode & 0o777) !== expectedMode) throw new Error(`${label} must be mode ${expectedMode.toString(8).padStart(4, '0')}`);
  if (expectedMode === null && (stat.mode & 0o444) === 0) throw new Error(`${label} must be readable`);
  return stat;
}
function json(file, label, expectedMode = 0o600) { regular(file, label, expectedMode); return JSON.parse(fs.readFileSync(file, 'utf8')); }
const state = json(statePath, 'deployment state');
const receipt = json(receiptPath, 'factory seed receipt');
const lineage = json(lineagePath, 'factory database lineage');
// ROM artifacts are read-only image inputs and are shipped mode 0644.  Only
// mutable factory evidence in /data is root-only mode 0600.
const provenance = json(provenancePath, 'factory provenance', null);
regular(seedPath, 'factory seed', null);
regular(publicationPath, 'factory database publication authority');
const receiptKeys = ['format', 'receiptKind', 'operationId', 'databasePath', 'seedSha256', 'dataMountIdentitySha256', 'databasePublicationSha256', 'createdDevice', 'createdInode', 'createdAt'];
if (!receipt || Object.keys(receipt).sort().join() !== receiptKeys.sort().join()
    || receipt.format !== 1 || receipt.receiptKind !== 'factory-seed' || receipt.databasePath !== databasePath
    || typeof receipt.operationId !== 'string' || receipt.operationId.length === 0) throw new Error('factory seed receipt is foreign or mismatched');
const lineageKeys = ['format', 'databasePath', 'seedReceiptSha256', 'seedSha256', 'dataMountIdentitySha256', 'databasePublicationSha256', 'createdDevice', 'createdInode', 'createdAt'];
if (!lineage || Object.keys(lineage).sort().join() !== lineageKeys.sort().join()
    || lineage.format !== 1 || lineage.databasePath !== databasePath) throw new Error('factory database lineage is foreign or mismatched');
const receiptSha = sha(fs.readFileSync(receiptPath));
const stateLibrary = require(stateLibraryPath);
const lineageSha = sha(Buffer.from(canonical(lineage)));
const provenanceSha = sha(fs.readFileSync(provenancePath));
stateLibrary.validateFactoryBaselineEnvelope(state, {
  expectedDeploymentId: receipt.operationId,
  expectedFactoryProvenanceSha256: provenanceSha,
  expectedSeedReceiptSha256: receiptSha,
  expectedDatabaseLineageSha256: lineageSha,
});
const stateLineage = state.parentDeployment.databaseLineage;
if (lineage.seedReceiptSha256 !== receiptSha || receipt.seedSha256 !== provenance.factorySeedSha256 || lineage.seedSha256 !== provenance.factorySeedSha256) throw new Error('factory seed receipt/lineage/provenance mismatch');
if (stateLineage.databaseLineageSha256 !== lineageSha || stateLineage.seedReceiptSha256 !== receiptSha) throw new Error('deployment state lineage does not bind receipt/lineage bytes');
if (sha(fs.readFileSync(seedPath)) !== provenance.factorySeedSha256) throw new Error('factory seed bytes do not match provenance');
const seed = require(seedLibraryPath);
seed.verifyFactoryDatabaseLineage(lineagePath, {
  database: databasePath,
  seedReceiptPath: receiptPath,
  expectedSeedReceiptSha256: receiptSha,
  expectedSeedSha256: provenance.factorySeedSha256,
  expectedDatabaseLineageSha256: lineageSha,
  databaseLineageState: stateLineage,
});
NODE

[ -e "$DB" ] || fail "database missing after verified 93 baseline"
if [ -f "$DB" ] && [ ! -L "$DB" ]; then
    exit 0
fi
fail "database is not a regular nonsymlink file"
