# conf/full_raspberrypi_bcm27xx_bcm2712/files/etc/uci-defaults/97_osi_db_seed
#!/bin/sh
set -eu

# The ROM seed helper owns hashing, SQLite checks, atomic publication, and the
# one-use receipt.  This UCI wrapper never overwrites a live database.
ROM_ROOT=${OSI_ROM_ROOT:-/rom}
DATA_ROOT=${OSI_DATA_ROOT:-/data}
DB="$DATA_ROOT/db/farming.db"
STATE_ROOT="$DATA_ROOT/osi-deploy"
SEED="$ROM_ROOT/usr/share/db/farming.db"
PROVENANCE="$ROM_ROOT/usr/share/osi-deploy/factory-image-provenance.json"
NODE=${OSI_NODE:-node}
PROFILE=${OSI_IMAGE_PROFILE:-$("$NODE" -e 'const fs=require("fs"); const p=JSON.parse(fs.readFileSync(process.argv[1])); process.stdout.write(p.profile);' "$PROVENANCE")}

[ -f "$SEED" ] || { logger -t osi-db-seed "ROM seed missing" 2>/dev/null || :; exit 1; }
[ -f "$PROVENANCE" ] || { logger -t osi-db-seed "ROM provenance missing" 2>/dev/null || :; exit 1; }
[ ! -e "$DB" ] || exit 0
for sidecar in "$DB-wal" "$DB-shm" "$DB-journal"; do
    [ ! -e "$sidecar" ] || { logger -t osi-db-seed "sidecar exists; refusing seed" 2>/dev/null || :; exit 1; }
done
mkdir -p "$DATA_ROOT/db" "$STATE_ROOT/receipts"
SEED_SHA256=$("$NODE" -e 'const fs=require("fs"); const p=JSON.parse(fs.readFileSync(process.argv[1])); process.stdout.write(p.factorySeedSha256);' "$PROVENANCE")
"$NODE" "$ROM_ROOT/usr/libexec/osi-factory-database-seed-cli.js" realize \
    --factory-seed "$SEED" --expected-seed-sha256 "$SEED_SHA256" --database "$DB" \
    --operation-id "factory-baseline-$PROFILE" \
    --receipt-out "$STATE_ROOT/receipts/factory-seed.json" \
    --database-lineage-out "$STATE_ROOT/database-lineage.json"
