# conf/full_raspberrypi_bcm27xx_bcm2712/files/etc/uci-defaults/99-default-node-red-flow
#!/bin/sh

DEFAULT_FLOW_SOURCE="/usr/share/flows.json"
NODE_RED_USER_DIR="/srv/node-red"
NODE_RED_FLOW_DEST="$NODE_RED_USER_DIR/flows.json"

mkdir -p "$NODE_RED_USER_DIR"

if [ ! -f "$NODE_RED_FLOW_DEST" ]; then
    cp "$DEFAULT_FLOW_SOURCE" "$NODE_RED_FLOW_DEST"
    chown node-red:node-red "$NODE_RED_FLOW_DEST"
    echo "Default Node-RED flow installed."
else
    echo "Node-RED flow already exists, not overwriting."
fi

exit 0 