# conf/full_raspberrypi_bcm27xx_bcm2712/files/etc/uci-defaults/99_set_sx1301_gateway_id
#!/bin/sh

[ -f /usr/libexec/osi-gateway-identity.sh ] && . /usr/libexec/osi-gateway-identity.sh

resolve_fallback_gateway_id() {
    local iface raw
    for iface in eth0 br-lan wlan0; do
        [ -r "/sys/class/net/$iface/address" ] || continue
        raw="$(cat "/sys/class/net/$iface/address" 2>/dev/null || true)"
        raw="$(normalize_gateway_eui "$raw" || true)"
        [ -n "$raw" ] && {
            printf '%s\n' "$raw"
            return 0
        }
    done
    return 1
}

gateway_identity_resolve >/dev/null 2>&1 || true
GATEWAY_ID="${GATEWAY_IDENTITY_DEVICE_EUI:-}"
[ -n "$GATEWAY_ID" ] || GATEWAY_ID="$(resolve_fallback_gateway_id || true)"
[ -n "$GATEWAY_ID" ] || exit 0

SECTION="$(gateway_identity_active_lora_section || true)"
[ -n "$SECTION" ] || exit 0

uci -q set "$SECTION.gateway_id=$GATEWAY_ID"
uci commit chirpstack-concentratord

exit 0
