# conf/base_raspberrypi_bcm27xx_bcm2708/files/etc/uci-defaults/99-initialize-custom-db
#!/bin/sh

CUSTOM_DB_SOURCE="/usr/share/sensor_data.db"
CUSTOM_DB_DEST="/srv/sensor_data.db"
SCHEMA_FILE="/usr/share/schema.sql"

mkdir -p /srv

if [ ! -f "$CUSTOM_DB_DEST" ]; then
    echo "Initializing custom sensor database..."
    cp "$CUSTOM_DB_SOURCE" "$CUSTOM_DB_DEST"
    sqlite3 "$CUSTOM_DB_DEST" < "$SCHEMA_FILE"
    # Optionally set permissions:
    chown node-red:node-red "$CUSTOM_DB_DEST"
    echo "Custom sensor database initialized."
else
    echo "Custom sensor database already exists."
fi

exit 0