# conf/full_raspberrypi_bcm27xx_bcm2712/files/etc/init.d/osi-lowpower
#!/bin/sh /etc/rc.common

START=99
STOP=10
USE_PROCD=1

SERVICE="/usr/libexec/osi-lowpower-window.sh"

lowpower_enabled() {
    [ "$(uci -q get osi-lowpower.main.enabled 2>/dev/null)" = "1" ]
}

start_service() {
    lowpower_enabled || return 0
    "$SERVICE" apply-wifi || true
    "$SERVICE" reconcile || true

    procd_open_instance
    procd_set_param command /bin/sh -c "while true; do $SERVICE reconcile; sleep 300; done"
    procd_set_param respawn 3600 5 0
    procd_close_instance
}

stop_service() {
    lowpower_enabled || return 0
    "$SERVICE" close || true
}

reload_service() {
    restart
}

status() {
    "$SERVICE" status
}

open() {
    "$SERVICE" open
}

close() {
    "$SERVICE" close
}
