# conf/full_raspberrypi_bcm27xx_bcm2709/files/etc/init.d/osi-bootstrap
#!/bin/sh /etc/rc.common

START=99

BOOTSTRAP_ROM="/usr/share/node-red/chirpstack-bootstrap.js"
BOOTSTRAP_FALLBACK="/srv/node-red/chirpstack-bootstrap.js"

stamp_valid() {
	[ -f /etc/osi-bootstrap.done ] || return 1
	[ -f /srv/node-red/.chirpstack.env ] || return 1
	grep -q 'CHIRPSTACK_APP_SENSORS=[0-9a-f]\{8\}-' /srv/node-red/.chirpstack.env 2>/dev/null || return 1
	return 0
}

boot() {
	start
}

start() {
	if stamp_valid; then
		return 0
	fi

	local bootstrap_script=""
	local provisioned_this_boot=0
	if [ -f "$BOOTSTRAP_ROM" ]; then
		bootstrap_script="$BOOTSTRAP_ROM"
	elif [ -f "$BOOTSTRAP_FALLBACK" ]; then
		bootstrap_script="$BOOTSTRAP_FALLBACK"
	fi

	if [ -z "$bootstrap_script" ]; then
		logger -t osi-bootstrap "bootstrap script not found, will retry next boot"
		return 0
	fi

	logger -t osi-bootstrap "waiting for ChirpStack gRPC..."
	local ready=0
	for i in $(seq 1 24); do
		if curl -sf --max-time 3 http://localhost:8080 2>/dev/null; then
			ready=1
			break
		fi
		sleep 5
	done

	if [ "$ready" = 0 ]; then
		logger -t osi-bootstrap "ChirpStack not ready after 120s, will retry next boot"
		return 0
	fi

	logger -t osi-bootstrap "running $(basename "$bootstrap_script")..."
	if node "$bootstrap_script"; then
		if touch /etc/osi-bootstrap.done; then
			provisioned_this_boot=1
			logger -t osi-bootstrap "bootstrap completed successfully"
		else
			logger -t osi-bootstrap "bootstrap completed but stamp write failed, will retry next boot"
		fi
	else
		logger -t osi-bootstrap "bootstrap failed, will retry next boot"
	fi

	if [ "$provisioned_this_boot" = 1 ]; then
		logger -t osi-bootstrap "requesting coordinated Node-RED restart after bootstrap"
		if ! /etc/init.d/osi-identityd ready; then
			rm -f /etc/osi-bootstrap.done
			logger -t osi-bootstrap "identityd is not ready; bootstrap will retry"
			return 0
		fi
		if /usr/libexec/osi-identityd.sh request-restart chirpstack_bootstrap 60; then
			logger -t osi-bootstrap "Node-RED restart requested after bootstrap"
		else
			rm -f /etc/osi-bootstrap.done
			logger -t osi-bootstrap "Node-RED restart request failed; bootstrap will retry"
			return 0
		fi
	fi
	return 0
}
