# files/etc/uci-defaults/99-node-red-setup
#!/bin/sh

# Install Node-RED and SQLite node
opkg update
opkg install node node-npm sqlite3-cli

# Install Node-RED globally
npm install -g node-red
npm install -g node-red-node-sqlite

# Create node-red user
adduser -D -h /var/lib/node-red -s /bin/false node-red

# Set up directories
mkdir -p /var/lib/node-red/.node-red
chown -R node-red:node-red /var/lib/node-red

# Copy default configuration if not exists
if [ ! -f /var/lib/node-red/.node-red/flows.json ]; then
    cp /usr/share/node-red/flows.json /var/lib/node-red/.node-red/ 2>/dev/null || true
    cp /usr/share/node-red/settings.js /var/lib/node-red/.node-red/ 2>/dev/null || true
    chown -R node-red:node-red /var/lib/node-red
fi

# Enable and start node-red service
/etc/init.d/node-red enable
/etc/init.d/node-red start

# Remove this script after first run
rm -f /etc/uci-defaults/99-node-red-setup

exit 0
EOF